# sv/sample/run
#!/bin/bash

sv start prestart && sleep 1 || exit 1

source /etc/envvars

exec 2>&1
exec sleep infinity
